# root/etc/s6-overlay/s6-rc.d/init-mastodon-config/up
/etc/s6-overlay/s6-rc.d/init-mastodon-config/run
