# root/etc/services.d/sidekiq/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

RAILS_ENV=production
PATH="${PATH}:/app/www/bin"

cd /app/www || exit 1

exec \
    s6-setuidgid abc /usr/bin/bundle exec sidekiq
